# misc/inittab
::sysinit:/etc/startup.sh
tty1::respawn:/bin/sh
::ctrlaltdel:/bin/busybox reboot
::shutdown:/bin/busybox swapoff -a
::shutdown:/bin/busybox umount -a -r
::restart:/etc/init
tty2::askfirst:/bin/sh
tty3::askfirst:/bin/sh

